# com.tedstudio.voiceinput.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tedstudio.voiceinput</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/.zshrc 2>/dev/null; cd /Users/tedstudio/Projects/whisper-voice-input; source venv/bin/activate; exec python voice_input.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/tedstudio/Projects/whisper-voice-input/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/tedstudio/Projects/whisper-voice-input/launchd.log</string>
</dict>
</plist>
